# .config/systemd/user/dwmstatus.service
[Unit]
Description=Updates dwm title

[Service]
Type=oneshot
ExecStart=/home/masonl/src/dwmstatus-rs/target/release/dwmstatus
