# .config/systemd/user/dwmstatus.timer
[Unit]
Description=Update dwm status

[Timer]
OnCalendar=*:*:0/30
